# pi/systemd/skygate-bypass.service
[Unit]
Description=SkyGate Aviation Bypass Daemon
Documentation=https://github.com/quartermint/skygate
After=network-online.target nftables.service
Wants=network-online.target
Requires=nftables.service

[Service]
Type=simple
ExecStart=/opt/skygate/skygate-bypass -config /data/skygate/bypass-domains.yaml -interval 60s
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skygate-bypass

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/skygate
PrivateTmp=true

[Install]
WantedBy=multi-user.target
