# pi/systemd/skygate-dashboard.service
[Unit]
Description=SkyGate Dashboard Daemon
Documentation=https://github.com/quartermint/skygate
After=network-online.target nftables.service pihole-FTL.service
Wants=network-online.target pihole-FTL.service
Requires=nftables.service

[Service]
Type=simple
ExecStart=/opt/skygate/skygate-dashboard -config /data/skygate/dashboard.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skygate-dashboard

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/skygate
PrivateTmp=true

[Install]
WantedBy=multi-user.target
